# jingui.service
[Unit]
Description=金匮要略临床思辨训练系统
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/JinGui/backend
Environment=OPENAI_BASE_URL=https://ai.linvk.com/v1
Environment=OPENAI_API_KEY=sk-zBM4DCxwpGMh58JaCTBWRFzPwvMIoGDaF1aW0JL7Ym8dpini
Environment=AI_MODEL=deepseek/deepseek-v4-flash
ExecStart=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
