# jingui.service
[Unit]
Description=金匮要略临床思辨训练系统
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/JinGui/backend
Environment=OPENAI_BASE_URL=https://api.deepseek.com
Environment=OPENAI_API_KEY=sk-aa2c5697a8214fcd909c6ed6e63f29bf
Environment=AI_MODEL=deepseek-chat
ExecStart=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
